# zfs-stats-exporter.service.in
[Unit]
Description=Export ZFS pool and dataset attributes for Prometheus

[Service]
Environment=PORT=9101
EnvironmentFile=-@SYSCONFDIR@/default/@NAME@
ExecStart=@SBINDIR@/@NAME@ $PORT

[Install]
WantedBy=multi-user.target
